# pi-music.service
[Unit]
Description=Pi Music Console – Touchscreen Music Player
After=multi-user.target sound.target
Wants=sound.target

[Service]
# Run as the auto-login user
User=pizza
Group=pizza

# Set the display so Tkinter can reach the framebuffer / X server
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pizza/.Xauthority
Environment=HOME=/home/pizza
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

WorkingDirectory=/home/pizza

# Use the kiosk launcher to handle both X and the Python app
ExecStart=/usr/bin/bash /home/pizza/Pi_Music_Console/start_kiosk.sh

# Restart automatically if it crashes
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-music

[Install]
WantedBy=multi-user.target
