# pi-music.service
[Unit]
Description=Pi Music Console – Touchscreen Music Player
After=multi-user.target sound.target
Wants=sound.target

[Service]
# Run as the auto-login user
User=pizza
Group=pizza

# Set the display so Tkinter can reach the framebuffer / X server
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pizza/.Xauthority
Environment=HOME=/home/pizza
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

WorkingDirectory=/home/pizza

# Wait until X is actually running
ExecStartPre=/bin/bash -c "until xset -q >/dev/null 2>&1; do sleep 1; done"

ExecStart=/usr/bin/python3 /home/pizza/Pi_Music_Console/music_player.py

# Restart automatically if it crashes
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-music

[Install]
WantedBy=multi-user.target
